# radserver0.socket
[Unit]
PartOf=radserver0.service

[Socket]
ListenFIFO=%t/radserver0.stdin