# xonotic.socket
[Unit]
PartOf=xonoitc.service

[Socket]
ListenFIFO=%t/xonotic.stdin